# app-pda/usbmuxd/files/usbmuxd
#!/sbin/runscript
# Copyright 1999-2010 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

depend() {
	need devfs
}

start() {
    ebegin "Starting usbmuxd"
    start-stop-daemon --start --exec /usr/sbin/usbmuxd -- -u -v -v
    eend $?
}

stop() {
    ebegin "Stopping usbmuxd"
    start-stop-daemon --stop --pidfile /var/run/usbmuxd.pid
    eend $?
}


restart() {
    stop
    start
}
